# com.user.transmission.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.user.transmission</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/transmission-daemon</string>
    <string>-f</string>
    <string>--config-dir</string>
    <string>/opt/homebrew/var/transmission/home</string>
  </array>

  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>

  <key>StandardOutPath</key><string>/tmp/transmission.out</string>
  <key>StandardErrorPath</key><string>/tmp/transmission.err</string>
</dict>
</plist>
